# tracker.timer
[Unit]
Description=Timer for Simple File Tracker

[Timer]
OnUnitActiveSec=5s
Unit=tracker.service

[Install]
WantedBy=timers.target
