# tracker.timer
[Unit]
Description=Timer for Simple File Tracker

[Timer]
OnBootSec=5s
OnUnitActiveSec=3s
Unit=tracker.service

[Install]
WantedBy=timers.target
